# ops/backup/systemd/hyperhash-wallet-backup.timer
[Unit]
Description=Run HyperHash wallet backup nightly

[Timer]
OnCalendar=*-*-* 02:30:00
RandomizedDelaySec=15m
Persistent=true
Unit=hyperhash-wallet-backup.service

[Install]
WantedBy=timers.target